# idlelocker.service
[Unit]
Description=Check if User has been Idle

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c '/usr/bin/idlelocker.sh'